# 3rd-party/prrte/contrib/docker/prte.service
[Unit]
Description=PRRTE persistent service.

[Service]
Type=simple
ExecStart=/opt/hpc/local/prrte/bin/prte

[Install]
WantedBy=multi-user.target
